# debian/postinst
#!/bin/sh

systemctl daemon-reload

#shellcheck source=./scripts/_bouncer.sh
. "/usr/lib/$DPKG_MAINTSCRIPT_PACKAGE/_bouncer.sh"
START=1

if [ "$1" = "configure" ]; then
    if need_api_key; then
        if ! set_api_key; then
            START=0
        fi
    fi
fi

systemctl --quiet is-enabled "$SERVICE" || systemctl unmask "$SERVICE" && systemctl enable "$SERVICE"

set_local_lapi_url 'CROWDSEC_LAPI_URL'

if ! getent passwd crowdsec-spoa >/dev/null; then
    adduser crowdsec-spoa --system --group --comment "crowdsec haproxy spoa bouncer"
fi

# Set config file group ownership and permissions
# Always set permissions to ensure file is readable by crowdsec-spoa group
# This is especially important on upgrades where old packages may have had incorrect permissions
if [ -f "$CONFIG" ]; then
    chmod 640 "$CONFIG" 2>/dev/null || true
    chgrp crowdsec-spoa "$CONFIG" 2>/dev/null || true
fi

if [ -d "/etc/haproxy" ]; then
    cp /usr/share/doc/crowdsec-haproxy-spoa-bouncer/examples/crowdsec.cfg /etc/haproxy/crowdsec.cfg
fi

# Selective restart logic for upgrades
# Only attempt restart for patch upgrades (same major.minor version)
STATE_DIR="/var/lib/$DPKG_MAINTSCRIPT_PACKAGE"
OLD_VERSION_FILE="${STATE_DIR}/old-version"
SHOULD_RESTART=0

if [ "$1" = "configure" ] && [ -n "${2:-}" ]; then
    # This is an upgrade (configure with old-version)
    OLD_VERSION="${2}"
    
    # If we have the old version from preinst, use it (more reliable)
    if [ -f "$OLD_VERSION_FILE" ]; then
        OLD_VERSION=$(cat "$OLD_VERSION_FILE")
        rm -f "$OLD_VERSION_FILE"
    fi
    
    # Get current version from dpkg
    CURRENT_VERSION=$(dpkg-query -W -f='${Version}' "$DPKG_MAINTSCRIPT_PACKAGE" 2>/dev/null || echo "")
    
    if [ -n "$OLD_VERSION" ] && [ -n "$CURRENT_VERSION" ]; then
        # Extract major.minor.patch from versions (strip debian revision if present)
        OLD_BASE=$(echo "$OLD_VERSION" | cut -d'-' -f1)
        CURRENT_BASE=$(echo "$CURRENT_VERSION" | cut -d'-' -f1)
        
        # Extract major.minor from both versions
        OLD_MAJOR_MINOR=$(echo "$OLD_BASE" | cut -d'.' -f1-2)
        CURRENT_MAJOR_MINOR=$(echo "$CURRENT_BASE" | cut -d'.' -f1-2)
        
        # If major.minor are the same, this is a patch upgrade
        if [ "$OLD_MAJOR_MINOR" = "$CURRENT_MAJOR_MINOR" ]; then
            SHOULD_RESTART=1
        fi
    fi
fi

# Attempt to restart service if it's a patch upgrade and service was running
if [ "$SHOULD_RESTART" -eq 1 ]; then
    # Check if service was active before upgrade (prerm would have stopped it)
    # We'll try to restart it regardless if it was enabled
    if systemctl --quiet is-enabled "$SERVICE" 2>/dev/null; then
        echo "Patch upgrade detected, attempting to restart service..."
        if systemctl restart "$SERVICE" 2>/dev/null; then
            echo "Service restarted successfully."
        else
            echo "⚠ Failed to restart service. Please start manually: systemctl start $SERVICE"
        fi
    fi
fi

echo ""
echo "=========================================="
echo "CrowdSec HAProxy SPOA Bouncer installed"
echo "=========================================="
echo ""

if [ "$START" -eq 0 ]; then
    echo "⚠ No API key was generated."
    echo "  Generate one with: cscli bouncers add <bouncer_name>"
    echo "  Add it to: $CONFIG"
    echo ""
fi

echo "Configuration: $CONFIG"
echo "Examples: /usr/share/doc/crowdsec-haproxy-spoa-bouncer/examples/"
echo "Documentation: https://docs.crowdsec.net/u/bouncers/haproxy_spoa"
echo ""
echo "Start bouncer: systemctl enable --now $SERVICE"
