# debian/postinst
#!/bin/sh

systemctl daemon-reload

#shellcheck source=./scripts/_bouncer.sh
. "/usr/lib/$DPKG_MAINTSCRIPT_PACKAGE/_bouncer.sh"
START=1

if [ "$1" = "configure" ]; then
    if need_api_key; then
        if ! set_api_key; then
            START=0
        fi
    fi
fi

systemctl --quiet is-enabled "$SERVICE" || systemctl unmask "$SERVICE" && systemctl enable "$SERVICE"

set_local_lapi_url 'CROWDSEC_LAPI_URL'

if ! getent passwd crowdsec-spoa >/dev/null; then
    adduser crowdsec-spoa --system --group --comment "crowdsec haproxy spoa bouncer"
fi

# Set config file group ownership
chgrp crowdsec-spoa "$CONFIG" 2>/dev/null || true

if [ -d "/etc/haproxy" ]; then
    cp /usr/share/doc/crowdsec-haproxy-spoa-bouncer/examples/crowdsec.cfg /etc/haproxy/crowdsec.cfg
fi

echo ""
echo "=========================================="
echo "CrowdSec HAProxy SPOA Bouncer installed"
echo "=========================================="
echo ""

if [ "$START" -eq 0 ]; then
    echo "⚠ No API key was generated."
    echo "  Generate one with: cscli bouncers add <bouncer_name>"
    echo "  Add it to: $CONFIG"
    echo ""
fi

echo "Configuration: $CONFIG"
echo "Examples: /usr/share/doc/crowdsec-haproxy-spoa-bouncer/examples/"
echo "Documentation: https://docs.crowdsec.net/u/bouncers/haproxy_spoa"
echo ""
echo "Start bouncer: systemctl enable --now $SERVICE"
