# config/crowdsec-spoa-bouncer-admin.socket
[Unit]
Description=CrowdSec SPOA Bouncer Admin Socket (Optional)
Documentation=https://docs.crowdsec.net/u/bouncers/haproxy_spoa/
# This socket is OPTIONAL and disabled by default
# To use it: systemctl enable --now crowdsec-spoa-bouncer-admin.socket
# Then restart the service: systemctl restart crowdsec-spoa-bouncer.service
Before=crowdsec-spoa-bouncer.service

[Socket]
ListenStream=/run/crowdsec-spoa/admin.sock
SocketMode=0600
SocketUser=root
SocketGroup=root
Service=crowdsec-spoa-bouncer.service
Accept=no

[Install]
WantedBy=sockets.target

